# experiment-setup/bootstrap/provisioning/clients/experiment_scenario_service
#!/sbin/openrc-run

description="Run the scenario Python script in a virtual environment"

depend() {
    after net            # Ensure the network is initialized
    after default          # Ensure all local scripts have finished
    after bootmisc       # Ensure miscellaneous boot tasks are complete
    need localmount      # Ensure all local filesystems are mounted
}

start() {
    ebegin "Starting scenario Python script"
    
    if [ -x /media/mmcblk0p2/init_experiment_client.sh ]; then
        /media/mmcblk0p2/init_experiment_client.sh
        eend $?
    else
        eend 1 "Script not found or not executable: /media/mmcblk0p2/init_experiment_client.sh"
    fi
}
